# rules/endpoint/linux/linux_cron_persistence.yml
schema_version: "1.0"
title: Linux cron / systemd unit installed by non-admin user
id: 7d8e9f0a-1b2c-4d3e-4f5a-6b7c8d9e0f1a
status: stable
level: medium
description: |
  Detects writes to common Linux persistence paths (/etc/cron.d,
  /etc/cron.hourly, /etc/cron.daily, /etc/systemd/system,
  ~/.config/systemd/user) by an executable not owned by root.
references:
  - https://attack.mitre.org/techniques/T1053/003/
  - https://attack.mitre.org/techniques/T1543/002/
tags:
  - attack.t1053.003
  - attack.t1543.002
  - attack.persistence
logsource:
  product: linux
  service: auditd
detection:
  selection:
    type: PATH
    name|startswith:
      - "/etc/cron"
      - "/var/spool/cron"
      - "/etc/systemd/system/"
      - "/home/"
    name|contains:
      - "systemd/user"
      - "cron"
  condition: selection
falsepositives:
  - Configuration management tools (Ansible, Puppet, Chef) running as a
    service account.
